# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=./gunicorn.sh

